# deploy/systemd/job-tracker-frontend.service
[Unit]
Description=Job Tracker frontend watcher container
Requires=docker.service
After=docker.service job-tracker-app.service

[Service]
Type=simple
WorkingDirectory=/opt/job-tracker
ExecStart=/usr/bin/env docker compose up --build --remove-orphans node
ExecStop=/usr/bin/env docker compose stop node
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
